# lib/init/kiln-controller.service
[Unit]
Description=kiln-controller

[Service]
Nice=-20
ExecStart=/home/pi/kiln-controller/venv/bin/python /home/pi/kiln-controller/kiln-controller.py

[Install]
WantedBy=multi-user.target
